# homelab/roles/truenas_nfs/templates/start-truenas-nfs.service.j2
[Unit]
Description=Start TrueNAS NFS Service via API
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/start-truenas-nfs.sh
TimeoutSec=30
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
