# homelab/roles/nfs-client/templates/nfs-mount.mount.j2
# {{ ansible_managed }}
# NFS mount unit for {{ nfs_mount.mount_point }}
# Generated by Ansible NFS client role

[Unit]
Description=NFS mount for {{ nfs_mount.server }}:{{ nfs_mount.path }}
After=network-online.target
Wants=network-online.target
Requires=rpcbind.service
After=rpcbind.service

[Mount]
What={{ nfs_mount.server }}:{{ nfs_mount.path }}
Where={{ nfs_mount.mount_point }}
Type={{ nfs_mount.fstype | default('nfs') }}
Options={{ nfs_mount.options | default(nfs_default_mount_options) }}
TimeoutSec=30

[Install]
WantedBy=multi-user.target
